# etc/systemd/system/rbd-target-gw.service.d/dependencies.conf
[Unit]
Wants=rbd-target-api.service
Requires=sys-kernel-config.mount
Before=rbd-target-api.service
