# etc/systemd/system/rbd-target-gw.service.d/dependencies.conf
[Unit]
Requires=sys-kernel-config.mount rbd-target-api.service
Before=rbd-target-api.service
